# deploy/hermes-tracy.service
# ABOUTME: Systemd unit managing hermes-tracy Docker container (Tracy instance).
# ABOUTME: Secrets decrypted from SOPS at start; stored in /run/secrets/ (tmpfs).

[Unit]
Description=Hermes Agent - Tracy (containerized)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
OnFailure=cron-alert@%n.service

[Service]
Type=oneshot
RemainAfterExit=true
User=root
Group=root
WorkingDirectory=/tank/services/active_services/hermes-agent/deploy

ExecStartPre=+/usr/local/bin/sops-decrypt.sh hermes-tracy /tank/services/active_services/hermes-tracy/.env.sops
ExecStart=/usr/bin/docker compose -f /tank/services/active_services/hermes-agent/deploy/docker-compose.yaml up -d --no-deps hermes-tracy
ExecStop=/usr/bin/docker compose -f /tank/services/active_services/hermes-agent/deploy/docker-compose.yaml stop hermes-tracy
ExecStopPost=+/bin/rm -rf /run/secrets/hermes-tracy

[Install]
WantedBy=multi-user.target
